# deploy/messenger-worker-prod.service
[Unit]
Description=Symfony Messenger Worker (production)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/cleanwhiskers.com
Environment=APP_ENV=prod
ExecStart=/usr/bin/php bin/console messenger:consume --time-limit=3600 --env=prod -vv
Restart=always
RestartSec=5s
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

